# root/etc/uci-defaults/80_accesscontrol
#!/bin/sh

touch /etc/config/example

return 0
